# docs/mitre-coverage-layer.json
{
  "name": "Cyber Incident Storyteller",
  "versions": {
    "layer": "4.4",
    "navigator": "4.9.5",
    "attack": "18"
  },
  "domain": "enterprise-attack",
  "description": "Techniques detectable from host logs \u2014 no cloud, no LLM, no SIEM required. Generated by: ais coverage --fmt navigator",
  "techniques": [
    {
      "techniqueID": "T1003",
      "score": 1,
      "comment": "cmd:mimikatz",
      "enabled": true
    },
    {
      "techniqueID": "T1003.001",
      "score": 1,
      "comment": "Windows Object Access, cmd:procdump",
      "enabled": true
    },
    {
      "techniqueID": "T1003.003",
      "score": 1,
      "comment": "cmd:ntdsutil",
      "enabled": true
    },
    {
      "techniqueID": "T1003.006",
      "score": 1,
      "comment": "Windows DS Object Access",
      "enabled": true
    },
    {
      "techniqueID": "T1003.008",
      "score": 1,
      "comment": "File Access",
      "enabled": true
    },
    {
      "techniqueID": "T1016",
      "score": 1,
      "comment": "cmd:ifconfig, cmd:ip, cmd:ipconfig",
      "enabled": true
    },
    {
      "techniqueID": "T1018",
      "score": 1,
      "comment": "cmd:arp, cmd:ping, cmd:nbtstat",
      "enabled": true
    },
    {
      "techniqueID": "T1021",
      "score": 1,
      "comment": "Windows Remote Logon, Windows Network Connection",
      "enabled": true
    },
    {
      "techniqueID": "T1021.002",
      "score": 1,
      "comment": "Windows Share Access",
      "enabled": true
    },
    {
      "techniqueID": "T1021.004",
      "score": 1,
      "comment": "Session Opened, cmd:ssh",
      "enabled": true
    },
    {
      "techniqueID": "T1025",
      "score": 1,
      "comment": "USB Connected",
      "enabled": true
    },
    {
      "techniqueID": "T1033",
      "score": 1,
      "comment": "cmd:whoami, cmd:id",
      "enabled": true
    },
    {
      "techniqueID": "T1046",
      "score": 1,
      "comment": "cmd:nmap, cmd:masscan",
      "enabled": true
    },
    {
      "techniqueID": "T1047",
      "score": 1,
      "comment": "cmd:wmic",
      "enabled": true
    },
    {
      "techniqueID": "T1048",
      "score": 1,
      "comment": "cmd:scp, cmd:rsync, cmd:ftp, cmd:sftp",
      "enabled": true
    },
    {
      "techniqueID": "T1049",
      "score": 1,
      "comment": "cmd:netstat, cmd:ss",
      "enabled": true
    },
    {
      "techniqueID": "T1053.002",
      "score": 1,
      "comment": "cmd:at",
      "enabled": true
    },
    {
      "techniqueID": "T1053.003",
      "score": 1,
      "comment": "Cron Execution, cmd:crontab",
      "enabled": true
    },
    {
      "techniqueID": "T1053.005",
      "score": 1,
      "comment": "Windows Scheduled Task, cmd:schtasks",
      "enabled": true
    },
    {
      "techniqueID": "T1057",
      "score": 1,
      "comment": "cmd:ps, cmd:tasklist",
      "enabled": true
    },
    {
      "techniqueID": "T1059",
      "score": 1,
      "comment": "cmd:nc, cmd:ncat, cmd:netcat, cmd:perl (+2 more)",
      "enabled": true
    },
    {
      "techniqueID": "T1059.001",
      "score": 1,
      "comment": "cmd:powershell, cmd:pwsh",
      "enabled": true
    },
    {
      "techniqueID": "T1059.003",
      "score": 1,
      "comment": "cmd:cmd",
      "enabled": true
    },
    {
      "techniqueID": "T1059.004",
      "score": 1,
      "comment": "Shell Execution, cmd:bash, cmd:sh",
      "enabled": true
    },
    {
      "techniqueID": "T1059.005",
      "score": 1,
      "comment": "cmd:wscript, cmd:cscript",
      "enabled": true
    },
    {
      "techniqueID": "T1059.006",
      "score": 1,
      "comment": "cmd:python, cmd:python3",
      "enabled": true
    },
    {
      "techniqueID": "T1069",
      "score": 1,
      "comment": "cmd:net, cmd:net1",
      "enabled": true
    },
    {
      "techniqueID": "T1070",
      "score": 1,
      "comment": "cmd:fsutil",
      "enabled": true
    },
    {
      "techniqueID": "T1070.001",
      "score": 1,
      "comment": "Windows Log Cleared",
      "enabled": true
    },
    {
      "techniqueID": "T1070.002",
      "score": 1,
      "comment": "cmd:shred, cmd:truncate",
      "enabled": true
    },
    {
      "techniqueID": "T1070.003",
      "score": 1,
      "comment": "cmd:history, cmd:unset",
      "enabled": true
    },
    {
      "techniqueID": "T1070.004",
      "score": 1,
      "comment": "File Deleted",
      "enabled": true
    },
    {
      "techniqueID": "T1071",
      "score": 1,
      "comment": "Network Connection, cmd:socat",
      "enabled": true
    },
    {
      "techniqueID": "T1078",
      "score": 1,
      "comment": "Accepted Password, Accepted Publickey, Audit Login, Admin Access (+1 more)",
      "enabled": true
    },
    {
      "techniqueID": "T1078.002",
      "score": 1,
      "comment": "Windows Privilege Assigned",
      "enabled": true
    },
    {
      "techniqueID": "T1082",
      "score": 1,
      "comment": "cmd:uname, cmd:hostname",
      "enabled": true
    },
    {
      "techniqueID": "T1083",
      "score": 1,
      "comment": "cmd:find, cmd:dir",
      "enabled": true
    },
    {
      "techniqueID": "T1087.002",
      "score": 1,
      "comment": "cmd:dsquery, cmd:dsget, cmd:adfind, cmd:ldifde (+1 more)",
      "enabled": true
    },
    {
      "techniqueID": "T1098",
      "score": 1,
      "comment": "Windows Group Member Added, Windows Account Changed, cmd:usermod",
      "enabled": true
    },
    {
      "techniqueID": "T1105",
      "score": 1,
      "comment": "cmd:wget, cmd:curl",
      "enabled": true
    },
    {
      "techniqueID": "T1110",
      "score": 1,
      "comment": "Failed Login, Auth Failure, Audit Auth Failure, Windows Kerberos PreAuth Failure (+1 more)",
      "enabled": true
    },
    {
      "techniqueID": "T1110.001",
      "score": 1,
      "comment": "Invalid User, Windows Logon Failure, cmd:hydra",
      "enabled": true
    },
    {
      "techniqueID": "T1110.002",
      "score": 1,
      "comment": "cmd:john, cmd:hashcat",
      "enabled": true
    },
    {
      "techniqueID": "T1112",
      "score": 1,
      "comment": "Windows Registry Modified, cmd:reg, cmd:regedit",
      "enabled": true
    },
    {
      "techniqueID": "T1115",
      "score": 1,
      "comment": "cmd:clip",
      "enabled": true
    },
    {
      "techniqueID": "T1124",
      "score": 1,
      "comment": "cmd:w32tm",
      "enabled": true
    },
    {
      "techniqueID": "T1132.001",
      "score": 1,
      "comment": "cmd:base64",
      "enabled": true
    },
    {
      "techniqueID": "T1134",
      "score": 1,
      "comment": "Windows Token Rights Adjusted",
      "enabled": true
    },
    {
      "techniqueID": "T1136.001",
      "score": 1,
      "comment": "Windows Account Created, cmd:useradd, cmd:adduser",
      "enabled": true
    },
    {
      "techniqueID": "T1140",
      "score": 1,
      "comment": "cmd:certutil",
      "enabled": true
    },
    {
      "techniqueID": "T1190",
      "score": 1,
      "comment": "Web Attack",
      "enabled": true
    },
    {
      "techniqueID": "T1197",
      "score": 1,
      "comment": "cmd:bitsadmin",
      "enabled": true
    },
    {
      "techniqueID": "T1201",
      "score": 1,
      "comment": "cmd:chage",
      "enabled": true
    },
    {
      "techniqueID": "T1218.001",
      "score": 1,
      "comment": "cmd:hh",
      "enabled": true
    },
    {
      "techniqueID": "T1218.003",
      "score": 1,
      "comment": "cmd:cmstp",
      "enabled": true
    },
    {
      "techniqueID": "T1218.004",
      "score": 1,
      "comment": "cmd:installutil",
      "enabled": true
    },
    {
      "techniqueID": "T1218.005",
      "score": 1,
      "comment": "cmd:mshta",
      "enabled": true
    },
    {
      "techniqueID": "T1218.007",
      "score": 1,
      "comment": "cmd:msiexec",
      "enabled": true
    },
    {
      "techniqueID": "T1218.008",
      "score": 1,
      "comment": "cmd:odbcconf",
      "enabled": true
    },
    {
      "techniqueID": "T1218.009",
      "score": 1,
      "comment": "cmd:regasm, cmd:regsvcs",
      "enabled": true
    },
    {
      "techniqueID": "T1218.010",
      "score": 1,
      "comment": "cmd:regsvr32",
      "enabled": true
    },
    {
      "techniqueID": "T1218.011",
      "score": 1,
      "comment": "cmd:rundll32",
      "enabled": true
    },
    {
      "techniqueID": "T1222",
      "score": 1,
      "comment": "cmd:chmod",
      "enabled": true
    },
    {
      "techniqueID": "T1482",
      "score": 1,
      "comment": "cmd:nltest",
      "enabled": true
    },
    {
      "techniqueID": "T1485",
      "score": 1,
      "comment": "cmd:cipher, cmd:sdelete",
      "enabled": true
    },
    {
      "techniqueID": "T1489",
      "score": 1,
      "comment": "Service Stopped",
      "enabled": true
    },
    {
      "techniqueID": "T1490",
      "score": 1,
      "comment": "cmd:vssadmin, cmd:bcdedit",
      "enabled": true
    },
    {
      "techniqueID": "T1505.003",
      "score": 1,
      "comment": "Web Shell",
      "enabled": true
    },
    {
      "techniqueID": "T1529",
      "score": 1,
      "comment": "cmd:shutdown",
      "enabled": true
    },
    {
      "techniqueID": "T1531",
      "score": 1,
      "comment": "Windows Account Deleted, cmd:passwd",
      "enabled": true
    },
    {
      "techniqueID": "T1543.002",
      "score": 1,
      "comment": "Service Started",
      "enabled": true
    },
    {
      "techniqueID": "T1543.003",
      "score": 1,
      "comment": "Windows Service Installed, cmd:sc",
      "enabled": true
    },
    {
      "techniqueID": "T1548.003",
      "score": 1,
      "comment": "Sudo Command",
      "enabled": true
    },
    {
      "techniqueID": "T1550.002",
      "score": 1,
      "comment": "Windows Explicit Credential Use, Windows NTLM Auth",
      "enabled": true
    },
    {
      "techniqueID": "T1552",
      "score": 1,
      "comment": "cmd:cmdkey",
      "enabled": true
    },
    {
      "techniqueID": "T1555",
      "score": 1,
      "comment": "cmd:vaultcmd",
      "enabled": true
    },
    {
      "techniqueID": "T1558",
      "score": 1,
      "comment": "Windows Kerberos TGT Request",
      "enabled": true
    },
    {
      "techniqueID": "T1558.003",
      "score": 1,
      "comment": "Windows Kerberos Service Ticket",
      "enabled": true
    },
    {
      "techniqueID": "T1560.001",
      "score": 1,
      "comment": "cmd:tar, cmd:zip, cmd:gzip",
      "enabled": true
    },
    {
      "techniqueID": "T1562.002",
      "score": 1,
      "comment": "cmd:auditpol, cmd:wevtutil",
      "enabled": true
    },
    {
      "techniqueID": "T1564.001",
      "score": 1,
      "comment": "cmd:attrib",
      "enabled": true
    },
    {
      "techniqueID": "T1595",
      "score": 1,
      "comment": "Web Scan",
      "enabled": true
    },
    {
      "techniqueID": "T1595.002",
      "score": 1,
      "comment": "Tool Fingerprint",
      "enabled": true
    }
  ],
  "gradient": {
    "colors": [
      "#ffffff",
      "#00b4d8"
    ],
    "minValue": 0,
    "maxValue": 1
  },
  "legendItems": [
    {
      "label": "Detected",
      "color": "#00b4d8"
    }
  ],
  "hideDisabled": false,
  "showTacticRowBackground": true,
  "tacticRowBackground": "#1a1a2e"
}
